# deploy/dockerfiles/api/rootfs/etc/s6-overlay/s6-rc.d/cron/run
#!/usr/bin/with-contenv

echo "* * * * * php /srv/app/artisan schedule:run" | crontab -u www-data -
exec crond -fl 2
